# etc/xdg/autostart/gnome-software-service.desktop
[Desktop Entry]
Type=Application
Name=GNOME Software
Exec=jhbuild run /opt/gnome/libexec/gnome-software-service
OnlyShowIn=GNOME;MATE;
